# deploy/luboss-frontend.service
[Unit]
Description=Luboss95 Village Banking v2 Frontend
After=network.target
Requires=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/luboss-vb/ui
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="NEXT_PUBLIC_BASE_PATH=/test"
Environment="NEXT_PUBLIC_API_URL=https://lubossvb.com/test/api"
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=luboss-frontend

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
